# images/guest/systemd/or3-bootstrap.service
[Unit]
Description=or3 guest bootstrap
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/or3-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
