# images/guest/systemd/or3-bootstrap.service
[Unit]
Description=or3 guest bootstrap
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/or3-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
